# hmi/barkeep/lib/barkeep.service
[Unit]
Description=Barkeep application for 'The Mewling Goat'
After=graphical.target

[Service]
ExecStart=/path/to/your/app
Restart=always
Environment=DISPLAY=:0

[Install]
WantedBy=graphical.target
